# server/services/monitoring-server.service
[Unit]
Description=ServerMon Monitoring Server
After=network.target

[Service]
Type=simple
User=monitoring
Group=monitoring
WorkingDirectory=/opt/monitoring-suhu/server
Environment="PORT=8080"
Environment="ENV=production"
Environment="DB_PATH=/var/lib/monitoring-suhu/monitoring.db"
ExecStart=/opt/monitoring-suhu/server/bin/monitoring-server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitoring-server

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/monitoring-suhu

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
